# supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
command=python /app/bot/main.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/app1.err.log
stdout_logfile=/var/log/app1.out.log

[program:scheduler]
command=python /app/scheduler.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/app2.err.log
stdout_logfile=/var/log/app2.out.log